# scripts/com.unified.ipsec.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.unified.ipsec.agent</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/unified-ipsec-agent/venv/bin/python</string>
        <string>-m</string>
        <string>agent.core</string>
        <string>config.json</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/opt/unified-ipsec-agent</string>
    <key>StandardOutPath</key>
    <string>/var/log/unified-ipsec-agent.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/unified-ipsec-agent.log</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
</dict>
</plist>
